# roles/bamboo/files/bamboo.ini
[program:bamboo]
redirect_stderr=true
command=/opt/bamboo/bamboo -bind=0.0.0.0:8000 -config=/opt/bamboo/config/production.json
